# provisioning/configuration/roles/buildbot-worker-macos/templates/buildbot.worker.plist.j2
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
         <!-- Adapted from: -->
         <!-- https://github.com/buildbot/buildbot-contrib/blob/4d9cdc6a6cc79abdcc0acb10ad545dd7bc2e840b/worker/contrib/os-x/net.sourceforge.buildbot.worker.plist -->
         <key>Label</key>
         <string>buildbot.worker</string>

         <key>UserName</key>
         <string>{{application_user}}</string>

         <key>WorkingDirectory</key>
         <string>/Users/{{application_user}}/worker</string>

         <!-- Extend PATH to include the location of worker scripts,
              `/usr/local/bin`
              Source: http://www.launchd.info/ -->
         <key>EnvironmentVariables</key>
         <dict>
                 <key>PATH</key>
                 <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
         </dict>

         <key>ProgramArguments</key>
         <array>
                 <string>/usr/local/bin/twistd</string>
                 <string>--nodaemon</string>
                 <string>--python=buildbot.tac</string>
                 <string>--logfile=buildbot.log</string>
                 <string>--prefix=worker</string>
         </array>

         <key>KeepAlive</key>
         <dict>
                 <key>SuccessfulExit</key>
                 <false/>
         </dict>

         <key>RunAtLoad</key>
         <true/>
</dict>
</plist>
