# provisioning/configuration/roles/buildbot-worker/templates/buildbot-worker.service.j2
# Adapted from:
# https://github.com/buildbot/buildbot-contrib/blob/a0f5d3898b9955c7d97b7dd2346a85e27252bb60/worker/contrib/systemd/buildbot-worker.service
[Unit]
Description=Buildbot Worker
Wants=network.target
After=network.target

[Service]
Type=forking
PIDFile={{home_dir}}/worker/twistd.pid
WorkingDirectory={{home_dir}}
ExecStartPre=/usr/local/bin/reboot-if-present.sh {{reboot_sentinel}}
ExecStart=/usr/bin/sudo -u {{application_user}} /usr/local/bin/buildbot-worker start worker
ExecReload=/usr/bin/sudo -u {{application_user}} /usr/local/bin/buildbot-worker restart worker
ExecStop=/usr/bin/sudo -u {{application_user}} /usr/local/bin/buildbot-worker stop worker
Restart=always
User=root

[Install]
WantedBy=multi-user.target
